# build/bootable/root/etc/systemd/system/appliance-tls.service
[Unit]
Description=Appliance TLS Certificate
Documentation=https://github.com/neutronbuild/neutron
After=appliance-mounts.target resizefs.service appliance-environment.service
Requires=appliance-mounts.target resizefs.service appliance-environment.service

[Service]
Type=oneshot
EnvironmentFile=/etc/vmware/environment
ExecStart=/etc/vmware/appliance-tls.sh

[Install]
WantedBy=appliance.target
