# build/bootable/root/etc/systemd/system/appliance-load-docker-images.service
[Unit]
Description=Load Cached Docker Images
After=appliance-mounts.target resizefs.service docker.service
Requires=appliance-mounts.target resizefs.service docker.service

[Service]
Type=oneshot
# allow other services to view this services as 'running' after exit
RemainAfterExit=yes
ExecStart=/etc/vmware/appliance-load-docker-images.sh
ExecStartPost=/usr/bin/systemctl disable appliance-load-docker-images.service

[Install]
WantedBy=appliance.target
