# build/bootable/root/etc/systemd/system/appliance-mounts.target
[Unit]
Description=Disk mounts
Documentation=https://github.com/neutronbuild/neutron
Requires=storage-data.mount storage-db.mount storage-log.mount repartition.service resizefs.service
After=storage-data.mount storage-db.mount storage-log.mount repartition.service resizefs.service

[Install]
WantedBy=appliance-ready.target local-fs.target
