# build/bootable/root/etc/systemd/system/appliance.target
[Unit]
Description=System customization for Appliance
Documentation=https://github.com/neutronbuild/neutron/blob/master/docs/DESIGN.md

Requires=multi-user.target
After=multi-user.target

AllowIsolate=yes
