# build/bootable/root/etc/systemd/system/appliance-network.service
[Unit]
Description=Network customization based on user provided parameter
Before=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
ExecStart=/etc/vmware/appliance-network.sh

[Install]
WantedBy=appliance-ready.target
